# deploy/mica-backup.timer
[Unit]
Description=Run the Mica backup periodically

[Timer]
# Daily at 03:00, with jitter; catch up if the machine was off at fire time.
OnCalendar=*-*-* 03:00:00
RandomizedDelaySec=300
Persistent=true

[Install]
WantedBy=timers.target
